# sv/supervisor/run
#!/bin/sh

sv start nimbus || exit 1

exec 2>&1
exec /storm/bin/storm supervisor
